# dashboard/dashboard-dashboard.service
# dashboard-dashboard.service – LEGACY template (kept for reference only).
#
# WARNING: This file is no longer installed by install-dashboard-service.sh.
#          It contains hard-coded placeholder values that will NOT work on
#          most installations.
#
# Use the installer instead, which auto-detects the correct user and paths:
#
#   sudo bash dashboard/install-dashboard-service.sh
#
# The installer writes mcbot-dashboard.service to /etc/systemd/system/ with
# the correct User, WorkingDirectory, and ExecStart values filled in.
#
# To uninstall:
#   sudo systemctl disable --now mcbot-dashboard.service
#   sudo rm /etc/systemd/system/mcbot-dashboard.service
#   sudo systemctl daemon-reload

[Unit]
Description=MCBOT Web Dashboard
# Start after the network is up (dashboard binds to 0.0.0.0)
After=network.target

[Service]
Type=simple
# Replace with the OS user that owns the repository (e.g. pi)
User=cyoa
# Replace with the absolute path to the cloned repository
WorkingDirectory=/home/cyoa/MCBOT
# Use the venv Python so all dependencies are available (run ./setup.sh to create .venv)
ExecStart=/home/cyoa/MCBOT/.venv/bin/python -m dashboard.app
# Restart the process automatically if it exits unexpectedly
Restart=on-failure
RestartSec=10
# Route all stdout/stderr to the system journal (view with journalctl)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
