# dashboard/mcbot-dashboard.service
# mcbot-dashboard.service – systemd unit file for the MCBOT web dashboard.
#
# NOTE: This file is a TEMPLATE with placeholder values.  Do NOT copy it
# directly to /etc/systemd/system/ – the paths will be wrong.
#
# Use the installer to generate a correctly-populated unit file:
#
#   bash dashboard/install-dashboard-service.sh
#
# Or run setup.sh (recommended) which calls the installer automatically.
#
# To uninstall:
#   sudo systemctl disable --now mcbot-dashboard.service
#   sudo rm /etc/systemd/system/mcbot-dashboard.service
#   sudo systemctl daemon-reload

[Unit]
Description=MCBOT Web Dashboard
# Require the network stack to be fully online so the dashboard is reachable
# from other devices on the local network immediately after boot.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Filled in at install time by install-dashboard-service.sh
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=__PYTHON__ -m dashboard.app
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
